# setup/cloudflare.service
[Unit]
Description=Cloudflare IP update
After=network.target 

[Service]
EnvironmentFile=/home/ubuntu/cloudflare.config
ExecStart=/bin/bash /home/ubuntu/mcbc-pco-app/setup/cloudflare.sh
Type=simple

[Install]
WantedBy=multi-user.target
